# deploy/telegram-mcp-backfill@.service
[Unit]
Description=Telegram MCP forward catch-up + bounded history backfill (%i) → Supabase
After=network-online.target telegram-mcp-%i.service
Wants=network-online.target

[Service]
Type=oneshot
User=idkras
WorkingDirectory=/home/idkras/telegram-mcp
EnvironmentFile=/etc/telegram-mcp/env.d/%i.env
Environment=PYTHONPATH=/home/idkras/telegram-mcp
ExecStart=/home/idkras/telegram-mcp/.venv/bin/python scripts/catch_up_recent_telegram_to_supabase.py --profile %i --limit-messages 5000 --deep-backfill-budget 1000 --deep-backfill-per-chat 250 --deep-backfill-select-limit 50
TimeoutStartSec=20min
Nice=10
