# scripts/com.rickai.telegram-deep-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ⚠️ DEPRECATED — НЕ загружать через launchctl bootstrap.

  Stage-7 RCA 2026-06-12 (B1 BLOCKING — backend-reviewer):
    Этот plist открывал ОТДЕЛЬНОЕ persistent TelegramClient соединение под
    той же session string параллельно с Mac MCP listener / laba persistent
    listener / 5-мин catch-up. Telegram считает каждое подключение по
    session string — два endpoints с одной строкой = AuthKeyDuplicated
    через несколько hourly прогонов. §Telegram session-per-endpoint
    invariant (RCA 2026-05-28) нарушен по построению.

  Новый канон (B1 fix):
    deep_backfill_all_chats() вызывается ВНУТРИ catch_up_recent cycle
    (см. heroes_platform/heroes_telegram_mcp/scripts/
    catch_up_recent_telegram_to_supabase.py --deep-backfill-budget=N).
    Один TelegramClient на cycle: start → forward catch-up → bounded
    backward deep-backfill (default 500 msgs/cycle) → disconnect.
    Соединение живёт секунды-минуту, параллельных listeners нет.

    launchd schedule: оставляем ТОЛЬКО com.rickai.telegram-periodic-catchup
    (5 мин). Этот файл — историческая ссылка + аварийный one-shot.

  Аварийный manual one-shot (НЕ launchd):
    Если нужен прицельный deep-backfill для clientID NNN (большая история,
    приоритет владельца) — запускайте ad-hoc:
        DEEP_BACKFILL_PRIORITY_CHATS=NNN \
        .venv/bin/python heroes_platform/heroes_telegram_mcp/scripts/\
        catch_up_recent_telegram_to_supabase.py \
        --profile ikrasinsky --deep-backfill-budget 5000 \
        --deep-backfill-per-chat 2000 --deep-backfill-select-limit 5
    Этот вызов делит соединение с listener-ом по одному endpoint (Mac),
    но НЕ открывает второй persistent — конфликта нет.

  Removal plan:
    1. На laba/проде убедиться что com.rickai.telegram-deep-backfill-*
       НЕ загружен (launchctl print gui/$UID/com.rickai.telegram-deep-backfill-*
       должно вернуть «could not find service»).
    2. Удалить ~/Library/LaunchAgents/com.rickai.telegram-deep-backfill-*.plist
       (если был установлен).
    3. Когда история всех клиентских чатов закроется backward-фазой
       periodic_catchup — этот файл можно удалить из репо.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rickai.telegram-deep-backfill-DEPRECATED</string>

  <key>Disabled</key>
  <true/>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/echo</string>
    <string>DEPRECATED: see file header. deep-backfill now runs inside periodic catch-up cycle.</string>
  </array>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
